# src/all-in-one/supervisor.d/api.ini
[program:api]
command=sh -c 'if [ "${SERVICE_API}" = "yes" ]; then exec /usr/share/bunkerweb/api/entrypoint.sh 2>&1 | sed "s/^/[API] /" | tee /var/log/bunkerweb/api.log; else echo "[API] Disabled by SERVICE_API setting" && exit 0; fi'
directory=/usr/share/bunkerweb/api
user=nginx
group=nginx
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/usr/share/bunkerweb/deps/python",LISTEN_ADDR="%(ENV_API_LISTEN_ADDR)s",LISTEN_PORT="%(ENV_API_LISTEN_PORT)s",FORWARDED_ALLOW_IPS="%(ENV_API_FORWARDED_ALLOW_IPS)s"
priority=25
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM
autostart=true
autorestart=false
redirect_stderr=false
killasgroup=true
stopasgroup=true
stdout_events_enabled=true
stderr_events_enabled=true
