# src/all-in-one/supervisor.d/autoconf.ini
[program:autoconf]
command=sh -c 'if [ "${AUTOCONF_MODE}" = "yes" ]; then exec /usr/share/bunkerweb/all-in-one/service-log-wrapper.sh autoconf /var/log/bunkerweb/autoconf.log "[AUTOCONF] " /usr/share/bunkerweb/autoconf/main.py; else echo "[AUTOCONF] Disabled by AUTOCONF_MODE setting" && exit 0; fi'
directory=/usr/share/bunkerweb/autoconf
user=nginx
group=nginx
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/usr/share/bunkerweb/deps/python"
priority=30
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM
autostart=true
autorestart=false
redirect_stderr=false
killasgroup=true
stopasgroup=true
stdout_events_enabled=true
stderr_events_enabled=true
