# src/all-in-one/supervisor.d/autoconf.ini
[program:autoconf]
command=sh -c 'if [ "${AUTOCONF_MODE}" = "yes" ]; then exec /usr/share/bunkerweb/autoconf/main.py 2>&1 | sed "s/^/[AUTOCONF] /"; else echo "[AUTOCONF] Disabled by AUTOCONF_MODE setting" && exit 0; fi'
directory=/usr/share/bunkerweb/autoconf
user=nginx
group=nginx
stdout_logfile=/var/log/bunkerweb/autoconf.stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/bunkerweb/autoconf.stderr.log
stderr_logfile_maxbytes=0
priority=30
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM
autostart=true
redirect_stderr=false
killasgroup=true
stopasgroup=true
stdout_events_enabled=true
stderr_events_enabled=true
