# src/all-in-one/supervisor.d/bunkerweb.ini
[program:bunkerweb]
command=sh -c 'exec /usr/share/bunkerweb/all-in-one/service-log-wrapper.sh bunkerweb /var/log/bunkerweb/bunkerweb.log "[BUNKERWEB] " /usr/share/bunkerweb/entrypoint.sh'
user=nginx
group=nginx
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/usr/share/bunkerweb/deps/python"
priority=10
startsecs=10
startretries=3
stopwaitsecs=30
stopsignal=TERM
autostart=true
autorestart=true
redirect_stderr=false
killasgroup=true
stopasgroup=true
stdout_events_enabled=true
stderr_events_enabled=true
