# src/all-in-one/supervisor.d/bunkerweb.ini
[program:bunkerweb]
command=sh -c 'exec /usr/share/bunkerweb/entrypoint.sh 2>&1 | sed "s/^/[BUNKERWEB] /"'
user=nginx
group=nginx
stdout_logfile=/var/log/bunkerweb/bunkerweb.stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/bunkerweb/bunkerweb.stderr.log
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/usr/share/bunkerweb/deps/python"
priority=10
startsecs=10
startretries=3
stopwaitsecs=30
stopsignal=TERM
autostart=true
autorestart=true
redirect_stderr=false
killasgroup=true
stopasgroup=true
stdout_events_enabled=true
stderr_events_enabled=true
