# src/all-in-one/supervisor.d/crowdsec.ini
[program:crowdsec]
command=sh -c 'if [ "${USE_CROWDSEC}" = "yes" ] && [[ "${CROWDSEC_API:-http://127.0.0.1:8000}" == http://127.0.0.1* || "${CROWDSEC_API:-http://127.0.0.1:8000}" == http://localhost* ]]; then exec /usr/share/bunkerweb/all-in-one/service-log-wrapper.sh crowdsec /var/log/bunkerweb/crowdsec.log "[CROWDSEC] " /usr/local/bin/crowdsec -c /etc/crowdsec/config.yaml; else echo "[CROWDSEC] Disabled by USE_CROWDSEC or CROWDSEC_API setting" && exit 0; fi'
user=nginx
group=nginx
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0
priority=12
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=TERM
autostart=true
autorestart=false
redirect_stderr=false
killasgroup=true
stopasgroup=true
stdout_events_enabled=true
stderr_events_enabled=true
