# src/all-in-one/supervisor.d/redis.ini
[program:redis]
command=sh -c 'if [ "${USE_REDIS}" = "yes" ] && { [ "${REDIS_HOST:-127.0.0.1}" = "127.0.0.1" ] || [ "${REDIS_HOST:-127.0.0.1}" = "localhost" ]; }; then exec /usr/share/bunkerweb/all-in-one/service-log-wrapper.sh redis /var/log/bunkerweb/redis.log "[REDIS] " redis-server /etc/redis.conf; else echo "[REDIS] Disabled by USE_REDIS or REDIS_HOST setting" && exit 0; fi'
user=nginx
group=nginx
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0
priority=12
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=TERM
autostart=true
autorestart=false
redirect_stderr=false
killasgroup=true
stopasgroup=true
stdout_events_enabled=true
stderr_events_enabled=true
