# src/all-in-one/supervisor.d/scheduler.ini
[program:scheduler]
command=sh -c 'if [ "${SERVICE_SCHEDULER}" = "yes" ]; then exec /usr/share/bunkerweb/all-in-one/service-log-wrapper.sh scheduler /var/log/bunkerweb/scheduler.log "[SCHEDULER] " /usr/share/bunkerweb/scheduler/entrypoint.sh; else echo "[SCHEDULER] Disabled by SERVICE_SCHEDULER setting" && exit 0; fi'
directory=/usr/share/bunkerweb/scheduler
user=nginx
group=nginx
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/usr/share/bunkerweb/deps/python:/usr/share/bunkerweb/db"
priority=30
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM
autostart=true
autorestart=false
redirect_stderr=false
killasgroup=true
stopasgroup=true
stdout_events_enabled=true
stderr_events_enabled=true
