# src/all-in-one/supervisor.d/ui.ini
[program:ui]
command=sh -c 'if [ "${SERVICE_UI}" = "yes" ]; then exec /usr/share/bunkerweb/ui/entrypoint.sh 2>&1 | sed "s/^/[UI] /"; else echo "[UI] Disabled by SERVICE_UI setting" && exit 0; fi'
directory=/usr/share/bunkerweb/ui
user=nginx
group=nginx
stdout_logfile=/var/log/bunkerweb/ui.stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/bunkerweb/ui.stderr.log
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/usr/share/bunkerweb/deps/python"
priority=20
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM
autostart=true
redirect_stderr=false
killasgroup=true
stopasgroup=true
stdout_events_enabled=true
stderr_events_enabled=true
