# src/all-in-one/supervisor.d/ui.ini
[program:ui]
command=sh -c 'if [ "${SERVICE_UI}" = "yes" ]; then exec /usr/share/bunkerweb/all-in-one/service-log-wrapper.sh ui /var/log/bunkerweb/ui.log "[UI] " /usr/share/bunkerweb/ui/entrypoint.sh; else echo "[UI] Disabled by SERVICE_UI setting" && exit 0; fi'
directory=/usr/share/bunkerweb/ui
user=nginx
group=nginx
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/usr/share/bunkerweb/deps/python"
priority=20
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM
autostart=true
autorestart=false
redirect_stderr=false
killasgroup=true
stopasgroup=true
stdout_events_enabled=true
stderr_events_enabled=true
