# src/all-in-one/supervisor.d/worker.ini
[program:worker]
command=sh -c 'if [ "${SERVICE_WORKER}" = "yes" ]; then exec /usr/share/bunkerweb/all-in-one/service-log-wrapper.sh worker "[WORKER] " /usr/share/bunkerweb/deps/python/bin/celery -A worker.app worker --loglevel=${LOG_LEVEL:-info} --concurrency=${WORKER_CONCURRENCY:-2} --pool=prefork --max-tasks-per-child=1 --max-memory-per-child=${WORKER_MAX_MEMORY_KB:-300000} --hostname=worker@%%h --without-heartbeat --without-mingle --without-gossip -Ofair -Q ${WORKER_QUEUES:-default,heavy}; else echo "[WORKER] Disabled by SERVICE_WORKER setting" && exit 0; fi'
directory=/usr/share/bunkerweb
user=nginx
group=nginx
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/1/fd/2
stderr_logfile_maxbytes=0
environment=PYTHONPATH="/usr/share/bunkerweb:/usr/share/bunkerweb/deps/python:/usr/share/bunkerweb/db:/usr/share/bunkerweb/utils:/usr/share/bunkerweb/api"
priority=28
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM
autostart=true
autorestart=false
redirect_stderr=false
killasgroup=true
stopasgroup=true
stdout_events_enabled=true
stderr_events_enabled=true
