# src/linux/rc.d/bunkerweb_scheduler
#!/bin/sh

# PROVIDE: bunkerweb_scheduler
# REQUIRE: LOGIN DAEMON NETWORKING
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable bunkerweb_scheduler:
#
# bunkerweb_scheduler_enable="YES"
#

. /etc/rc.subr

name="bunkerweb_scheduler"
rcvar="${name}_enable"
desc="BunkerWeb Scheduler service"

load_rc_config $name

: ${bunkerweb_scheduler_enable:="NO"}
: ${bunkerweb_scheduler_user:="root"}
: ${bunkerweb_scheduler_group:="nginx"}
: ${bunkerweb_scheduler_pidfile:="/var/run/bunkerweb/scheduler.pid"}
: ${bunkerweb_scheduler_wrapper_pidfile:="/var/run/bunkerweb/scheduler-wrapper.pid"}

pidfile="${bunkerweb_scheduler_pidfile}"
start_cmd="${name}_start"
stop_cmd="${name}_stop"
reload_cmd="${name}_reload"
status_cmd="${name}_status"
extra_commands="reload"

bunkerweb_scheduler_start()
{
    echo "Starting ${name}."
    mkdir -p /var/run/bunkerweb
    /usr/sbin/daemon -f -p "${bunkerweb_scheduler_wrapper_pidfile}" /usr/local/bin/bash /usr/share/bunkerweb/scripts/bunkerweb-scheduler.sh start
}

bunkerweb_scheduler_stop()
{
    if [ -f "${pidfile}" ]; then
        echo "Stopping ${name}."
        /usr/local/bin/bash /usr/share/bunkerweb/scripts/bunkerweb-scheduler.sh stop
    else
        echo "${name} is not running."
    fi
    if [ -f "${bunkerweb_scheduler_wrapper_pidfile}" ]; then
        wrapper_pid=$(cat "${bunkerweb_scheduler_wrapper_pidfile}" 2>/dev/null)
        [ -n "${wrapper_pid}" ] && kill "${wrapper_pid}" 2>/dev/null || true
        rm -f "${bunkerweb_scheduler_wrapper_pidfile}"
    fi
}

bunkerweb_scheduler_reload()
{
    echo "Reloading ${name}."
    /usr/local/bin/bash /usr/share/bunkerweb/scripts/bunkerweb-scheduler.sh reload
}

bunkerweb_scheduler_status()
{
    if [ -f "${pidfile}" ]; then
        pid=$(cat "${pidfile}")
        if kill -0 "${pid}" 2>/dev/null; then
            echo "${name} is running as pid ${pid}."
            return 0
        fi
    fi
    echo "${name} is not running."
    return 1
}

run_rc_command "$1"
